# code/xquakejs/bin/proxy.service
[Unit]
Description=Quake proxy service.
[Service]
Type=simple
ExecStart=node /home/ioq3srv/Quake3e/quakejs/bin/proxy.js
Restart=on-failure
[Install]
WantedBy=multi-user.target
